# profile.d/path.sh
#!/bin/sh

export PATH=$PATH:$HOME/.local/share/bin:$HOME/.local/bin:/var/lib/flatpak/exports/bin:$HOME/.config/emacs/bin
